# elias-web.service
[Unit]
Description=Elias Web Console
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/elias-web
ExecStart=/usr/bin/node node_modules/tsx/dist/cli.mjs src/server.ts
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
